# systemctl/daily-backup.service
[Unit]
Description=Daily Backup For Vitor

[Service]
EnvironmentFile=/home/vitor/.config/daily-backup/env
ExecStart=/usr/bin/daily-backup launch